# scripts/maria.service
# M.A.R.I.A. AI Agent - V3 Unified Launcher
# Runs daemon (homeostasis tick loop) + Web UI in one process.
# Replaces separate maria.service + maria-ui.service.
#
# Copy to /etc/systemd/system/maria.service
# Then: sudo systemctl daemon-reload && sudo systemctl enable maria
# Disable old: sudo systemctl disable maria-ui

[Unit]
Description=M.A.R.I.A. AI Agent (V3)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=maria
WorkingDirectory=/home/maria/maria
ExecStart=/home/maria/maria/venv/bin/python maria.py
Restart=on-failure
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30

EnvironmentFile=/home/maria/maria/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=maria

[Install]
WantedBy=multi-user.target
