# scripts/maria.service
# M.A.R.I.A. - systemd unit (runs the daemon + Web UI in one process)
#
# This is a TEMPLATE. Before enabling, edit the User and the three absolute
# paths below to match your system, then:
#   sudo cp scripts/maria.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now maria

[Unit]
Description=M.A.R.I.A. AI Agent
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple

# EDIT: the Linux user Maria should run as
User=YOUR_USERNAME

# EDIT: absolute path to your clone (all three lines below)
WorkingDirectory=/absolute/path/to/Maria
ExecStart=/absolute/path/to/Maria/venv/bin/python maria.py
EnvironmentFile=/absolute/path/to/Maria/.env

Restart=on-failure
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=maria

[Install]
WantedBy=multi-user.target
