# scripts/maria-ui.service
# OBSOLETE - Web UI is now integrated into maria.py (V3 Unified Launcher)
# To disable: sudo systemctl disable maria-ui
# To remove: sudo rm /etc/systemd/system/maria-ui.service && sudo systemctl daemon-reload
#
# If you need Web UI ONLY (without daemon), use:
#   python maria.py --ui
#
# This file is kept for reference. Do not enable alongside maria.service.

[Unit]
Description=M.A.R.I.A. Web UI (OBSOLETE - use maria.service instead)
After=network.target maria.service
Wants=maria.service

[Service]
Type=simple
User=maria
WorkingDirectory=/home/maria/maria
ExecStart=/home/maria/maria/venv/bin/python run_ui.py
Restart=on-failure
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30

EnvironmentFile=/home/maria/maria/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=maria-ui

[Install]
WantedBy=multi-user.target
